# vws.service.in
[Unit]
Description=Virtual Workstation Autostart
Documentation=man:vws(1)
After=network.target dnsmasq.service
[Service]
Type=forking
ExecStart=@bindir@/vws autostart
ExecStop=@bindir@/vws shutdown --wait
Restart=no
[Install]
WantedBy=multi-user.target
Alias=vws.service
